# reminders/mcu-monitor@.service
[Unit]
Description=MCU monitor (%I)
After=dev-%i.device
BindsTo=dev-%i.device

[Service]
Type=simple
User=robot
Group=robot
SupplementaryGroups=robotdata dialout
UMask=0002
EnvironmentFile=/etc/robot_path.conf
ExecStart=/opt/robot/scripts/mcu_monitor.sh %i
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target

