# reminders/update-repo.service
[Unit]
Description=update a repository under ROBOT_PATH
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=robot
EnvironmentFile=/etc/robot_path.conf
ExecStart=/home/robot/update_repo.sh /etc/robot_repos.conf
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
